# 99-xrandr-trigger.rules
# send screenwatch events on add/remove
ACTION=="change", SUBSYSTEM=="drm", RUN+="screenwatch-event"
